# actividad4/saludo-serv.service
[Unit]
Description=Saludo Service 

[Service]
ExecStart=/usr/local/bin/Script.sh

[Install]
WantedBy=multi-user.target